# files/udev-rules-90-batbarupdate.rules
SUBSYSTEM=="power_supply", RUN+="/usr/bin/pkill -RTMIN+2 i3blocks"
